# vm/rootfs/etc/init.d/containerfy-compose
#!/sbin/openrc-run

description="Load preloaded images and start compose services"

depend() {
    need docker
    after docker
}

start() {
    ebegin "Starting compose services"

    # Load preloaded images on first boot
    if [ -d /var/cache/containerfy/images ] && [ ! -f /data/.images-loaded ]; then
        einfo "Loading preloaded container images..."
        for tar in /var/cache/containerfy/images/*.tar; do
            [ -f "$tar" ] || continue
            einfo "  Loading $(basename "$tar")..."
            docker load -i "$tar" 2>/dev/null
        done
        mkdir -p /data
        touch /data/.images-loaded
    fi

    if [ -f /etc/containerfy/docker-compose.yml ]; then
        docker compose -f /etc/containerfy/docker-compose.yml up -d
    fi

    eend $?
}

stop() {
    ebegin "Stopping compose services"
    if [ -f /etc/containerfy/docker-compose.yml ]; then
        docker compose -f /etc/containerfy/docker-compose.yml down --timeout 15
    fi
    eend $?
}
